# deploy/rpg.service
# Systemd unit file for Rpg daemon mode.
#
# Installation:
#   1. Copy this file to /etc/systemd/system/rpg.service
#   2. Copy rpg binary to /usr/local/bin/rpg
#   3. Create config: /etc/rpg/config.toml
#   4. Create user:   useradd --system --no-create-home rpg
#   5. Enable:        systemctl enable --now rpg
#
# Environment variables (set in /etc/rpg/env or override file):
#   PGHOST, PGPORT, PGUSER, PGDATABASE, PGPASSWORD (or .pgpass)
#   ANTHROPIC_API_KEY (optional, for AI features)

[Unit]
Description=Rpg — self-driving Postgres agent
Documentation=https://github.com/NikolayS/project-alpha
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=rpg
Group=rpg

ExecStart=/usr/local/bin/rpg --daemon --config /etc/rpg/config.toml
ExecReload=/bin/kill -HUP $MAINPID

# Environment
EnvironmentFile=-/etc/rpg/env

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

# Allow writing PID file and logs
ReadWritePaths=/run/rpg /var/log/rpg

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitIntervalSec=60

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpg

[Install]
WantedBy=multi-user.target
